# deploy/systemd/genie-ai-runtime.service
[Unit]
Description=GeniePod AI Runtime (Jetson-tuned LLM, OpenAI-compatible)
Documentation=https://github.com/GeniePod/genie-ai-runtime
After=network.target
# Claim the LLM KV cache before memory-heavy voice/container services start.
# Jetson testing for issue #75 showed the same `-c 4096` request fitting only
# ~1.7k ctx after the full stack was resident, but fitting 4k/6k/8k ctx when
# genie-ai-runtime loaded first.
Before=genie-whisper.service genie-whisper-warmup.service homeassistant.service genie-core.service
ConditionPathExists=/opt/geniepod/bin/jetson-llm-server
# Conflicts with genie-llm.service: both bind :8080. systemd will refuse
# to start the second one while the first is running, so a misconfigured
# [services.llm] block fails loudly instead of silently double-binding.
Conflicts=genie-llm.service

[Service]
Type=simple
# Keep the GGUF in page cache across restarts when the kernel can. Clearing
# VM caches here made every runtime restart cold-load Qwen3 again (issue #69).
# Use INT8 KV so the Jetson service can reserve an 8k context on Orin Nano
# when systemd starts it before memory-heavy services. `-c` is still clamped
# by runtime memory budget, so boot/start ordering matters.
ExecStart=/opt/geniepod/bin/jetson-llm-server \
    -m ${GENIEPOD_LLM_MODEL} \
    -p 8080 \
    -c ${GENIEPOD_AI_RUNTIME_CONTEXT} \
    --int8-kv
Environment=GENIEPOD_LLM_MODEL=/opt/geniepod/models/Qwen3-4B-Q4_K_M.gguf
Environment=GENIEPOD_AI_RUNTIME_CONTEXT=8192
Restart=on-failure
RestartSec=5
TimeoutStartSec=120

# GPU needs full system access
ProtectSystem=no
SupplementaryGroups=video render

[Install]
WantedBy=geniepod.target
